# alpine-virtualmail/resources/etc/services.d/05_crond/run
#!/usr/bin/execlineb -P

exec /usr/sbin/crond -f -S
