# alpine-openfire/resources/etc/services.d/50-openfire/run
#!/command/with-contenv sh

exec /app/openfire/bin/openfire run
