# alpine-fuglu/resources/etc/services.d/00_syslog/run
#!/command/execlineb -P

exec syslogd -n -S -O -
