# alpine-fuglu/resources/etc/services.d/00_syslog/run
#!/usr/bin/execlineb -P

exec syslogd -n -S -O -
