# alpine-spiderfoot/resources/etc/services.d/10-spiderfeet/run
#!/usr/bin/with-contenv sh

cd /app

exec /usr/bin/python /app/sf.py 0.0.0.0:5001
