# alpine-teamspeak3/resources/etc/services.d/10_teamspeak3/run
#!/usr/bin/execlineb -P

/app/ts3server_minimal_runscript.sh
