# alpine-teamspeak3/resources/etc/services.d/10_teamspeak3/run
#!/usr/bin/with-contenv sh

exec /bin/s6-setuidgid teamspeak3 /app/ts3server_minimal_runscript.sh inifile=ts3server.ini
