# alpine-grimd/resources/etc/services.d/10_grimd/run
#!/command/with-contenv sh

cd /config

su-exec grimd /usr/local/bin/grimd -config /config/grimd.toml -update
