# alpine-virtualmail/resources/etc/services.d/10_dovecot/run
#!/command/execlineb -P

exec /usr/sbin/dovecot -c /etc/dovecot/dovecot.conf -F
