# alpine-virtualmail/resources/etc/services.d/00_syslog/run
#!/command/execlineb -P

exec /sbin/syslogd -n -S -O -
