# alpine-virtualmail/resources/etc/services.d/00_syslog/run
#!/usr/bin/with-contenv sh

exec syslogd -n -S
