# alpine-nginx-php7/resources/etc/services.d/php-fpm7/run
#!/command/with-contenv sh

exec /usr/sbin/php-fpm7
