# alpine-nginx-phpbb-php7/resources/etc/crontabs/root
*/5 * * * * test ! -d /data/web/install && /usr/bin/phpbbcli cron:run > /dev/stdout 2>&1
