# alpine-bitlbee/resources/etc/services.d/bitlbee/run
#!/usr/bin/with-contenv sh

exec bitlbee -u nobody -Dvn
