# alpine-posty-webui/resources/etc/services.d/10_posty/run
#!/command/with-contenv sh

cd /data/webui

test -f /data/config/lighttpd.user && exec lighttpd -f /etc/lighttpd/lighttpd-htpasswd.conf -D

exec /usr/sbin/lighttpd -f /etc/lighttpd/lighttpd.conf -D
