# alpine-nginx-phpbb/resources/etc/services.d/crond/run
#!/command/with-contenv sh

exec /usr/sbin/crond -f -S -c /etc/crontabs
