# alpine-virtualmail-fuglu/resources/etc/services.d/60_postfix/run
#!/usr/bin/execlineb -P

exec /usr/libexec/postfix/master -d
