# alpine-mariadb/resources/etc/services.d/mariadb/run
#!/usr/bin/with-contenv sh

exec /usr/bin/mysqld_safe --pid-file=$PID_FILE
