# alpine-weecloud/resources/etc/services.d/weecloud/run
#!/usr/bin/with-contenv sh

cd /app/weecloud

exec /usr/bin/node app.js
