# alpine-weechat/resources/etc/services.d/10-dropbear/run
#!/command/with-contenv sh

exec /usr/sbin/dropbear -wsjkBFER
