# alpine-weechat/resources/etc/services.d/10-dropbear/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/dropbear -wsjkBFER
