# alpine-haproxy/resources/etc/services.d/10_syslog/run
#!/usr/bin/execlineb -P

/sbin/syslogd -n -O -
