# alpine-haproxy/resources/etc/services.d/10_syslog/run
#!/command/execlineb -P

/sbin/syslogd -n -O -
