# alpine-virtualmail-fuglu/resources/etc/crontabs/vmail
0 8,20  *  *  * /usr/local/bin/learnSA.sh
