# alpine-softether/resources/etc/services.d/10-vpnserver/run
#!/command/with-contenv sh

LANG=en_US.UTF-8

exec /usr/local/bin/vpnserver execsvc
