# alpine-nginx-phpbb/resources/etc/crontabs/root
*/5 * * * * s6-setuidgid www-data php -f /data/web/bin/phpbbcli.php cron:run > /dev/stdout 2>&1
