# alpine-fluentd/resources/etc/services.d/10-fluentd/run
#!/command/with-contenv sh

exec /command/s6-setuidgid fluent /usr/bin/fluentd -c /fluent/fluent.conf -v
