# alpine-fluentd/resources/etc/services.d/10-fluentd/run
#!/usr/bin/with-contenv sh

exec /bin/s6-setuidgid fluent /usr/bin/fluentd -c /fluent/fluent.conf -v
