# alpine-mumble/resources/etc/services.d/10_mumble/run
#!/command/with-contenv sh

exec /usr/bin/murmurd -fg -v
