# alpine-mumble/resources/etc/services.d/10_mumble/run
#!/usr/bin/with-contenv sh

exec /usr/bin/murmurd -fg -v
