# alpine-fuglu/resources/etc/services.d/20_spamd/run
#!/command/execlineb -P

exec /usr/sbin/spamd --username vmail --groupname vmail --create-prefs --max-children 5 --helper-home-dir --allow-tell
