# alpine-fuglu/resources/etc/services.d/20_spamd/run
#!/usr/bin/execlineb -P

exec /usr/sbin/spamd --username vmail --groupname vmail
