# alpine-netdata/resources/etc/services.d/10_netdata/run
#!/command/with-contenv sh

exec /usr/sbin/netdata -D -u root -s /host -i 0.0.0.0 -p 19999