# alpine-netdata/resources/etc/services.d/10_netdata/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/netdata -D -u root -s /host -p 19999 
