# alpine-virtualmail/resources/etc/services.d/25_clamd/run
#!/command/execlineb -P

exec /usr/sbin/clamd
