# alpine-nginx-php7/resources/etc/services.d/nginx/run
#!/command/with-contenv sh

exec /usr/sbin/nginx
