# alpine-prosody/resources/etc/services.d/50-prosody/run
#!/command/with-contenv sh

/sbin/su-exec prosody /usr/bin/prosody
