# alpine-prosody/resources/etc/services.d/50-prosody/run
#!/usr/bin/with-contenv sh

su-exec prosody /usr/bin/prosody
