# alpine-posty-api/resources/etc/services.d/10_posty/run
#!/command/with-contenv sh

cd /data/api

exec bundle exec rackup
