# alpine-nodejs-weecloud/resources/etc/services.d/weecloud/run
#!/usr/bin/with-contenv sh

cd /app

exec /usr/bin/node app.js
