# alpine-aptly/resources/etc/services.d/50-nginx/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/nginx -g 'daemon off; pid /tmp/nginx.pid;'
