# alpine-nginx-freshrss/resources/etc/crontabs/root
*/5 * * * * /command/s6-setuidgid www-data php -f /data/web/app/actualize_script.php > /dev/stdout 2>&1
