# alpine-nginx-nextcloud/resources/etc/services.d/cron/run
#!/usr/bin/with-contenv sh

/usr/sbin/crond -f -S -l 1 -c /etc/crontabs
