# alpine-aptly/resources/etc/services.d/10-api/run
#!/command/with-contenv sh

exec /command/s6-setuidgid aptly /app/aptly/aptly -config /etc/aptly.conf api serve
