# alpine-aptly/resources/etc/services.d/10-api/run
#!/usr/bin/with-contenv sh

exec /app/aptly/aptly api serve
