# alpine-aptly/resources/etc/services.d/10-api/run
#!/usr/bin/with-contenv sh

exec /bin/s6-setuidgid aptly /app/aptly/aptly -config /etc/aptly.conf api serve
