# alpine-haproxy/resources/etc/services.d/50_haproxy/run
#!/usr/bin/with-contenv sh

exec haproxy -W -db -p /tmp/haproxy.pid -f /etc/haproxy/haproxy.cfg -f /data/haproxy.d/
