# alpine-haproxy/resources/etc/services.d/50_haproxy/run
#!/usr/bin/with-contenv sh

OPT="-f /etc/haproxy/haproxy.cfg -f /data/haproxy.d/"

exec haproxy -db -p /tmp/haproxy.pid $OPT
