# alpine-djigzo/resources/etc/services.d/10-backend/run
#!/command/with-contenv sh

exec /app/djigzo/wrapper/wrapper-djigzo -c ../../djigzo.wrapper.conf $@
