# alpine-weechat/resources/etc/services.d/20-tor/run
#!/command/with-contenv sh

su-exec weechat /usr/bin/tor
