# alpine-weechat/resources/etc/services.d/20-tor/run
#!/usr/bin/with-contenv sh

exec /usr/bin/tor
