# alpine-weechat/resources/etc/services.d/20-tor/run
#!/usr/bin/with-contenv sh

su-exec weechat /usr/bin/tor
