# alpine-weechat/resources/etc/services.d/20-tor/run
#!/command/with-contenv sh

/sbin/su-exec weechat /usr/bin/tor
