# alpine-bigbrother-b3/resources/etc/services.d/b3/run
#!/command/with-contenv sh

HOME=/home/b3 exec /command/s6-setuidgid b3 /usr/bin/python2 /b3/b3_run.py --config /data/b3.ini
