# alpine-nuclide-remote/resources/etc/services.d/10_sshd/run
#!/usr/bin/with-contenv sh

if [ ! -d /etc/dropbear ]; then
	mkdir /etc/dropbear
fi

exec dropbear -F -R -E -B
