# alpine-musicbot/resources/etc/services.d/10_musicbot/run
#!/command/with-contenv sh

cd /app/MusicBot

exec /command/s6-setuidgid musicbot /usr/bin/python3 /app/MusicBot/run.py
