# alpine-musicbot/resources/etc/services.d/10_musicbot/run
#!/usr/bin/with-contenv sh

cd /app/MusicBot

exec /bin/s6-setuidgid musicbot /usr/bin/python3 /app/MusicBot/run.py
