# alpine-nginx-nextcloud/resources/etc/crontabs/www-data
*/15  *  *  *  * php -f /data/web/cron.php
