# alpine-haproxy/resources/etc/services.d/99_docker_gen/run
#!/command/execlineb -P

with-contenv
fdmove -c 2 1

multisubstitute
{
  importas -u -D "unix:///tmp/docker.sock" DOCKER_ENDPOINT DOCKER_ENDPOINT
}
/app/docker-gen -config /app/docker-gen.cfg -endpoint "${DOCKER_ENDPOINT}"
