# alpine-virtualmail/resources/etc/services.d/30_assp/run
#!/usr/bin/with-contenv sh

cd /usr/share/assp/

exec perl assp.pl
