# alpine-virtualmail/resources/etc/services.d/30_assp/run
#!/usr/bin/execlineb -P

cd /usr/share/assp/

exec perl assp.pl
