# alpine-virtualmail/resources/etc/services.d/30_assp/run
#!/command/execlineb -P

cd /usr/share/assp/

exec /usr/bin/perl assp.pl
