# alpine-fuglu/resources/etc/services.d/50_fuglu/run
#!/command/with-contenv sh
exec /usr/bin/fuglu --foreground --config /data/fuglu/fuglu.conf
