# alpine-fuglu/resources/etc/services.d/50_fuglu/run
#!/usr/bin/with-contenv sh
exec /usr/bin/fuglu --foreground --config=/etc/fuglu/fuglu.conf
