# alpine-spiderfoot/resources/etc/services.d/10-spiderfoot/run
#!/command/with-contenv sh

cd /app

exec /command/s6-setuidgid spiderfoot /usr/bin/python3 /app/sf.py -l 0.0.0.0:5001
